# Systemd Service Unit Files/pigpiod.service
[Unit]
Description=pigpio daemon starten

[Service]
Type=forking
ExecStart=/usr/bin/pigpiod

[Install]
WantedBy=multi-user.target
